# wipsync-daemon.service
[Unit]
Description=wipsync — background dirty-repo checker

[Service]
Type=simple
ExecStart=%h/.local/bin/wipsync-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
